# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:spacewalk-setup]
command=/opt/spacewalk.sh
autorestart=false
